# systemd/system/archrepo@.path
[Path]
PathChanged=/srv/http/archrepo/%i

[Install]
WantedBy=multi-user.target
